# etc/systemd/system/lxc@android.service.d/override.conf
[Service]
TimeoutStopSec=5s
ExecStop=
ExecStop=/usr/bin/lxc-stop -k -n android
